# etc/e11-cli/e11/e11admin/attack_student_instances.cron
0 * * * *  cd /Users/simsong/gits/csci_e-11/staff/lab1 && HIDDEN=csci-e-11-veritas-E11 venv/bin/python attack_student_instances.py
